# examples/handsfree/hf-hs@.service
[Unit]
Description=Bluetooth Handsfree Audio Service %I

[Service]
Type=simple
Environment=ARGS=%I
ExecStart=/usr/local/bin/handsfree.bash $ARGS

